# deploy/yks-backend.service
[Unit]
Description=YKS Takip Backend
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/yks-tracker
WorkingDirectory=/var/www
Environment=DB_PATH=/var/db/uygulama.db
Environment=FRONTEND_PATH=/var/www/yks-frontend
Environment=TELEGRAM_BOT_TOKEN=
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
